# cron/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
PERL5LIB=/usr/local/lib/perl

# Run DB cleanup task every night
0 5 * * * root . /etc/cron_env; /usr/local/bin/clean_logs.pl >> /var/log/cron.log 2>&1
